# files/etc/systemd/system/frr-ready.path
[Unit]
Description=Monitor for fabric ip readyness

[Path]
PathExists=/run/fabric_ip_ready
Unit=frr.service

[Install]
WantedBy=network.target